# etc/teamcity.target.dist
[Unit]
Description=Teamcity Server + agents
Requires=teamcity-agent.target
Requires=teamcity-server.service
